# configs/vpn.cron
*/1 * * * * bash -x /opt/vpn/vpncronrun
